# deploy/sms_pro_celery_worker.service.example
[Unit]
Description=SMS Pro Celery Worker
After=network.target mysql.service redis-server.service
Requires=mysql.service redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/sms_pro
EnvironmentFile=/var/www/sms_pro/.env
ExecStart=/var/www/sms_pro/venv/bin/celery -A sms_platform worker --loglevel=INFO --concurrency=4
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
